# recipes_flask.conf
description "Starts Flask as a Tornado IOLoop service"
author "Fred"
start on vagrant-ready
stop on runlevel [06]
exec python /vagrant/recipes_tornado.py
respawn