# recipes_rethinkdb.conf
description "Starts RethinkDB"
author "Fred"
start on vagrant-ready
stop on runlevel [06]
exec rethinkdb --bind all
respawn